# packaging/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /opt/vidispine-content-viewer/uwsgi.ini --plugin python --die-on-term
environment = DJANGO_SETTINGS_MODULE=app.settings.prod,HTTPS=on,wsgi.url_scheme=https
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command = /usr/sbin/nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
